# scripts/launchd/com.jarvis.kos-patrol.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jarvis.kos-patrol</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/chenyuanquan/.bun/bin/bun</string>
        <string>run</string>
        <string>/Users/chenyuanquan/Projects/jarvis-knowledge-os-v2/skills/kos-jarvis/kos-patrol/run.ts</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/chenyuanquan/.bun/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
        <key>GOOGLE_GENERATIVE_AI_API_KEY</key>
        <string><FILL:NANO_BANANA_API_KEY></string>
        <key>GBRAIN_EMBEDDING_MODEL</key>
        <string>google:gemini-embedding-001</string>
        <key>GBRAIN_EMBEDDING_DIMENSIONS</key>
        <string>1536</string>
    </dict>
    <!-- Daily 08:07 (off the :00 mark to avoid thundering herd) -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>7</integer>
    </dict>
    <key>RunAtLoad</key>
    <false/>
    <!-- patrol's exit code is by design: 0=clean, 1=ERROR (real fail), 2=WARN-only.
         Treat 0+2 as success so launchd doesn't flag normal warn-only runs. -->
    <key>SuccessfulExitCodes</key>
    <array>
        <integer>0</integer>
        <integer>2</integer>
    </array>
    <key>StandardOutPath</key>
    <string>/Users/chenyuanquan/Projects/jarvis-knowledge-os-v2/skills/kos-jarvis/kos-patrol/patrol.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/chenyuanquan/Projects/jarvis-knowledge-os-v2/skills/kos-jarvis/kos-patrol/patrol.stderr.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/chenyuanquan/Projects/jarvis-knowledge-os-v2</string>
</dict>
</plist>
