# scripts/launchd/com.jarvis.enrich-sweep.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jarvis.enrich-sweep</string>
    <!--
      Weekly enrich-sweep: scans the whole brain, extracts entities via Haiku,
      and creates Tier-2/Tier-3 stubs. Keys live in .env.local (sourced inside
      the wrap script). Tight --min-mentions 3 default matches what produced
      a clean wave on 2026-04-17.

      Wrapped via Minions shell-job for the 1h hard timeout + audit trail.
      Direct bun would have no safety net if the NER loop runs away (account
      pool exhaustion, infinite retry, etc). Wrap also pins the exact arg
      shape (--min-mentions 3 --max-tier2 30) inside the script, not the
      plist, so re-tuning doesn't require a launchctl bootout/bootstrap.

      Manual invocation:
        bash scripts/minions-wrap/enrich-sweep.sh
        bun run skills/kos-jarvis/enrich-sweep/run.ts --plan --min-mentions 3
    -->
    <key>ProgramArguments</key>
    <array>
        <string>/Users/chenyuanquan/Projects/jarvis-knowledge-os-v2/scripts/minions-wrap/enrich-sweep.sh</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/chenyuanquan/.bun/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
        <key>GOOGLE_GENERATIVE_AI_API_KEY</key>
        <string><FILL:NANO_BANANA_API_KEY></string>
        <key>GBRAIN_EMBEDDING_MODEL</key>
        <string>google:gemini-embedding-001</string>
        <key>GBRAIN_EMBEDDING_DIMENSIONS</key>
        <string>1536</string>
    </dict>
    <!-- Sunday 22:13 local. After end-of-week reading, before bedtime. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>22</integer>
        <key>Minute</key>
        <integer>13</integer>
    </dict>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/Users/chenyuanquan/Projects/jarvis-knowledge-os-v2/skills/kos-jarvis/enrich-sweep/sweep.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/chenyuanquan/Projects/jarvis-knowledge-os-v2/skills/kos-jarvis/enrich-sweep/sweep.stderr.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/chenyuanquan/Projects/jarvis-knowledge-os-v2</string>
</dict>
</plist>
