# docs/guides/minions-deployment-snippets/systemd.service
[Unit]
Description=Cortex worker
Documentation=https://github.com/Versatly/CortexBrain/blob/main/docs/guides/minions-deployment.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=cortex
Group=cortex
WorkingDirectory=/srv/cortex
EnvironmentFile=/etc/cortex.env
ExecStart=/usr/local/bin/cortex jobs supervisor --concurrency 2
Restart=always
RestartSec=10s
KillSignal=SIGTERM
TimeoutStopSec=30s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cortex-worker
LimitNOFILE=65535
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/srv/cortex /home/cortex/.cortex

[Install]
WantedBy=multi-user.target
