# scripts/launchd/_archived/com.jarvis.kos-compat-api.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jarvis.kos-compat-api</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/chenyuanquan/.bun/bin/bun</string>
        <string>run</string>
        <string>/Users/chenyuanquan/Projects/jarvis-knowledge-os-v2/server/kos-compat-api.ts</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/chenyuanquan/.bun/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
        <key>KOS_API_PORT</key>
        <string>7225</string>
        <key>KOS_API_TOKEN</key>
        <string>ed855ebd6989e0b286ddafa24094ef0a054fec9b71ab774a97d3339d58172c35</string>
        <key>GBRAIN_HOME</key>
        <string>/Users/chenyuanquan/brain</string>
        <key>OPENAI_BASE_URL</key>
        <string>http://127.0.0.1:7222/v1</string>
        <key>OPENAI_API_KEY</key>
        <string>stub-for-gemini-shim</string>
        <key>ANTHROPIC_BASE_URL</key>
        <string>https://crs.chenge.ink/api</string>
    </dict>
    <key>KeepAlive</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>StandardOutPath</key>
    <string>/Users/chenyuanquan/Projects/jarvis-knowledge-os-v2/server/kos-compat-api.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/chenyuanquan/Projects/jarvis-knowledge-os-v2/server/kos-compat-api.stderr.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/chenyuanquan/Projects/jarvis-knowledge-os-v2</string>
</dict>
</plist>
