# scripts/launchd/_archived/com.jarvis.kos-compat-api.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jarvis.kos-compat-api</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/chenyuanquan/.bun/bin/bun</string>
        <string>run</string>
        <string>/Users/chenyuanquan/Projects/jarvis-knowledge-os-v2/server/kos-compat-api.ts</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/chenyuanquan/.bun/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
        <key>KOS_API_PORT</key>
        <string>7225</string>
        <key>KOS_API_TOKEN</key>
        <string><FILL:KOS_API_TOKEN></string>
        <key>GOOGLE_GENERATIVE_AI_API_KEY</key>
        <string><FILL:NANO_BANANA_API_KEY></string>
        <key>GBRAIN_EMBEDDING_MODEL</key>
        <string>google:gemini-embedding-001</string>
        <key>GBRAIN_EMBEDDING_DIMENSIONS</key>
        <string>1536</string>
        <key>ANTHROPIC_BASE_URL</key>
        <string>https://crs.chenge.ink/api</string>
    </dict>
    <key>KeepAlive</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>StandardOutPath</key>
    <string>/Users/chenyuanquan/Projects/jarvis-knowledge-os-v2/server/kos-compat-api.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/chenyuanquan/Projects/jarvis-knowledge-os-v2/server/kos-compat-api.stderr.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/chenyuanquan/Projects/jarvis-knowledge-os-v2</string>
</dict>
</plist>
